# .config/systemd/user/podman-auto-update.service.d/override.conf
[Service]
Restart=on-failure
RestartSec=1s
RestartSteps=10
RestartMaxDelaySec=30min
