# .config/systemd/user/waybar.service.d/override.conf
[Unit]
After=wayland-session@sway.desktop.target

[Service]
Slice=app-graphical.slice

[Install]
WantedBy=
WantedBy=wayland-session@sway.desktop.target
