# .etcbuild/systemd-flatpak-update.timer
[Unit]
Description=Update Flatpak

[Timer]
OnCalendar=weekly
AccuracySec=1d
RandomizedDelaySec=1d
Persistent=true

[Install]
WantedBy=timers.target
