# .config/containers/systemd/languagetool.container
# vim: ft=systemd

[Unit]
Description=Containerized LanguageTool HTTP server

[Container]
Image=docker.io/erikvl87/languagetool:latest
AutoUpdate=registry

DropCapability=all
NoNewPrivileges=true
ReadOnly=true

Environment=Java_Xms=512m
Environment=Java_Xmx=1g

PublishPort=0.0.0.0:8081:8010/tcp

PodmanArgs=--cpus=2

[Service]
Restart=on-failure

[Install]
WantedBy=default.target
