# .config/containers/systemd/languagetool.container
[Unit]
Description=Containerized LanguageTool HTTP server

[Container]
Image=docker.io/erikvl87/languagetool:latest
AutoUpdate=registry
PublishPort=[::1]:8081:8010/tcp

[Service]
Restart=on-failure

[Install]
WantedBy=default.target
