# .config/containers/systemd/languagetool.container
# vim: ft=systemd

[Unit]
Description=Containerized LanguageTool HTTP server

[Container]
Image=docker.io/meyay/languagetool:latest
AutoUpdate=registry

DropCapability=all
NoNewPrivileges=true
UserNS=keep-id
ReadOnly=true

Tmpfs=/tmp
Volume=%D/fasttext:/fasttext

PublishPort=0.0.0.0:8081:8081/tcp

PodmanArgs=--cpus=1

[Service]
ExecStartPre=/usr/bin/mkdir -p %D/fasttext
Restart=on-failure

[Install]
WantedBy=default.target
