# .config/containers/systemd/languagetool.container
[Unit]
Description=Containerized LanguageTool HTTP server

[Container]
Image=docker.io/erikvl87/languagetool:latest
PublishPort=127.0.0.1:8081:8010/tcp
PublishPort=[::1]:8081:8010/tcp

[Install]
WantedBy=default.target
